# ops/systemd/thril-printer-guard.service
[Unit]
Description=THRIL printer default guard
After=network-online.target cups.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/thril-printer-guard
ExecStart=/usr/local/bin/thril-printer-guard.sh

[Install]
WantedBy=multi-user.target
